# Library/Application Support/BackupMinder/org.backupminder.sample.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.watchmanmonitoring.backupminder.sample</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/BackupMinder/BackupMinder</string>
		<string>--Name</string>
		<string>Config Name</string>
		<string>--BackupSource</string>
		<string>/Path/To/Backups</string>
		<string>--ArchiveDestination</string>
		<string>/Path/To/Archive</string>
		<string>--NameContains</string>
		<string>ImportantFile</string>
		<string>--BackupsToLeave</string>
		<string>30</string>
		<string>--WarnDays</string>
		<string>7</string>
	</array>
	<key>WatchPaths</key>
	<array>
		<string>/Users/gary/BackupDir</string>
	</array>
</dict>
</plist>
